# .eqmod/autopilot/com.eqmod.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
EQMOD autopilot watchdog launchd job.

Runs once per hour. Checks last_tick.txt and queue state; mails Michael
once per day with a summary and immediately on a hard problem (no tick for 12h+,
HUMAN_NEEDED count grew, etc.).

Install:
    cp .eqmod/autopilot/com.eqmod.watchdog.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.eqmod.watchdog.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.eqmod.watchdog</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/mkupermann/GitHub/EQMOD/.eqmod/autopilot/run_watchdog.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>KeepAlive</key>
    <false/>
    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/mkupermann/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/mkupermann</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/mkupermann/GitHub/EQMOD</string>

    <key>StandardOutPath</key>
    <string>/Users/mkupermann/.eqmod/autopilot/watchdog.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/mkupermann/.eqmod/autopilot/watchdog.err</string>
</dict>
</plist>
